# flagscd.service
[Unit]
Description=A Flags client daemon

[Service]
ExecStart=/usr/bin/flagscd
User=flagscd

[Install]
WantedBy=network.target
